# ops/systemd/mold-ledger.timer
[Unit]
Description=Daily Mold ledger drop (09:00 UTC)

[Timer]
OnCalendar=*-*-* 09:00:00 UTC
Persistent=true
RandomizedDelaySec=300

[Install]
WantedBy=timers.target
